# k8s/weatherflow-udp-forwarder.service
[Unit]
Description=WeatherFlow UDP Packet Forwarder
Documentation=https://github.com/dgorman/weatherflow-dashboards-aio
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/weatherflow
ExecStart=/usr/bin/python3 -u /opt/weatherflow/udp_interface_forwarder.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=false
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
